# cron/crontab
#crontab

*/2 * * * * user tool_cron/get_ext_rates
*/10 * * * * user tool_cron/proc_history